# packages/090_init/files/unitfiles/mdevd-coldplug.service
[Unit]
Description= Initialize kernel modules
After=mdevd.service

[Service]
ExecStart=/usr/local/sbin/mdevd-coldplug
Type=oneshot

[Install]
WantedBy=onboot.target
