# packages/090_init/files/unitfiles/mdevd.service
[Unit]
Description= Start mdevd daemon

[Service]
ExecStart=/usr/local/sbin/mdevd
Type=simple
Restart=always

[Install]
WantedBy=onboot.target
