# packages/090_init/files/unitfiles/onboot.target
[Unit]
Description= Onboot services checkpoint

[Install]
WantedBy=default.target
